# etc/mythpublish.cron
# select, publish then inform once every 3 hours
15  */3 *   *   *   /home/mythtv/bin/mythpublish_select
