# etc/nadtmyth_wakeup_times.upstart
description "set mythtv wakeup times"

start on started mysql

task
exec /usr/bin/sudo -u mythtv -i /home/mythtv/bin/nadtmyth_wakeup_times
